# ft_services/srcs/influxdb/data/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:influx]
command=sh -c "kill $(pidof influx) || influxd --config /etc/influxdb.conf && kill $(pidof supervisord)"
autostart=true
autorestart=true
priority=1

[program:telegraf]
command=sh -c "kill $(pidof telegraf) || /usr/bin/telegraf --config /etc/telegraf/telegraf.conf && kill $(pidof supervisord)"
autostart=true
autorestart=true
priority=2

[program:database]
command=sh -c "sleep 3 | echo CREATE DATABASE telegraf | influx"
autostart=true
autorestart=unexpected
exitcodes=0
priority=3