# ft_services/srcs/phpmyadmin/data/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm7]
command=sh -c "kill $(pidof php-fpm7) || /usr/sbin/php-fpm7 -R -F && kill $(pidof supervisord)"
autostart=true
autorestart=true

[program:nginx]
command=sh -c "kill $(pidof nginx) || nginx -g 'daemon off;' && kill $(pidof supervisord)"
autostart=true
autorestart=true